# scripts/suspend-systemd/suspend.timer
# Put inside /etc/systemd/system
[Unit]
Description=Suspend if user is idle

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
